# ansible/templates/fastlanes.service.j2
[Unit]
Description=FastLanes Next.js app (standalone)
After=network.target

[Service]
Type=simple
User={{ app_user }}
Group={{ app_user }}
# current — симлинк на активный релиз. После switch-симлинка
# рестарт службы поднимет уже новый код.
WorkingDirectory={{ app_dir }}/current
EnvironmentFile={{ app_dir }}/current/.env
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
# Стандартные env'ы Next.js standalone-сервера.
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
Environment=NODE_ENV=production

# Безопасность
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
